# overlay/etc/systemd/system/contest-clean-home.service
[Unit]
Description=Contest Runtime - Reset user home directory to clean state
DefaultDependencies=no
After=local-fs.target
Before=display-manager.service graphical.target
ConditionKernelCommandLine=contest.clean_home

[Service]
Type=oneshot
ExecStart=/usr/lib/contest/clean-home.sh
StandardOutput=journal+console
StandardError=journal+console
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
