# overlay/etc/systemd/system/contest-uninstall.service
[Unit]
Description=Contest Runtime - Remove deployment from local disk
# Must complete before the display manager starts; the service reboots
# when done so the user never reaches the desktop.
DefaultDependencies=no
After=local-fs.target
Before=display-manager.service graphical.target
ConditionKernelCommandLine=contest.uninstall

[Service]
Type=oneshot
ExecStart=/usr/lib/contest/uninstall.sh
StandardOutput=journal+console
StandardError=journal+console
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
