# overlay/etc/systemd/system/stats-report.service
[Unit]
Description=Stats Report Sender

[Service]
Type=oneshot
ExecStart=/usr/local/bin/stats-send.sh
